# hack/aws/userdata/control_plane/files/etc/systemd/system/kubelet.service.d/10-cgroup.conf
# /etc/systemd/system/kubelet.service.d/10-cgroup.conf
# Sets the cgroup for the kubelet service
[Service]
CPUAccounting=true
MemoryAccounting=true
Slice=podruntime.slice
