# hack/aws/userdata/workers/files/etc/systemd/system/disable-swap.service
[Unit]
After=local-fs.target

[Service]
ExecStart=/sbin/swapoff -a

[Install]
WantedBy=multi-user.target
